# sec-ops-05-endpoint-security/final-quiz.xml
<?xml version="1.0" encoding="UTF-8"?>
<quiz>
  <title>Final Quiz: Endpoint Security Technologies</title>
  <session>05</session>
  <type>final</type>
  <time_minutes>15</time_minutes>
  <description>Assess mastery of Session 05 content on Endpoint Security Technologies</description>

  <question id="1">
    <text>A SOC analyst sees Windows Event ID 4624 with LogonType 3 and AuthenticationPackageName NTLM for a Domain Administrator account, originating from an external IP address at 2 AM. Which combination of factors makes this MOST suspicious?</text>
    <options>
      <option id="A" correct="false">LogonType 3 is always suspicious</option>
      <option id="B" correct="true">NTLM authentication for a domain account at an unusual hour from an external IP</option>
      <option id="C" correct="false">The event occurred at 2 AM</option>
      <option id="D" correct="false">Domain Administrator accounts always generate Event 4624</option>
    </options>
    <explanation>The combination is suspicious because: (1) Domain accounts should use Kerberos, not NTLM — NTLM for a domain account can indicate pass-the-hash; (2) external IP means the source is outside the corporate network; (3) 2 AM is outside normal business hours. Any one factor alone might be explainable, but all three together is a strong indicator of unauthorized access.</explanation>
    <difficulty>intermediate</difficulty>
    <topic>Windows Event Log Analysis</topic>
    <mitre_technique>T1078 — Valid Accounts, T1550.002 — Pass the Hash</mitre_technique>
  </question>

  <question id="2">
    <text>An analyst investigates a Sysmon Event 10 (ProcessAccess) alert. The source process is C:\Users\john\AppData\Local\Temp\svc.exe accessing C:\Windows\System32\lsass.exe with GrantedAccess 0x1FFFFF. What is MOST likely happening?</text>
    <options>
      <option id="A" correct="false">A legitimate antivirus product is scanning for infections</option>
      <option id="B" correct="false">Windows is performing a routine credential cache refresh</option>
      <option id="C" correct="true">A credential dumping tool (e.g., Mimikatz) is attempting to extract password hashes</option>
      <option id="D" correct="false">LSASS is terminating abnormally due to a crash</option>
    </options>
    <explanation>Three indicators confirm this is credential dumping: (1) The source process is in AppData/Temp — a user-writable path used for malware staging; (2) GrantedAccess 0x1FFFFF is PROCESS_ALL_ACCESS — full access to the LSASS process memory; (3) LSASS is the credential store. Legitimate security products (AV) do have specific LSASS access patterns but are excluded by vendor in most EDR configurations. An unsigned binary from Temp with PROCESS_ALL_ACCESS is definitively malicious.</explanation>
    <difficulty>intermediate</difficulty>
    <topic>EDR Telemetry Analysis</topic>
    <mitre_technique>T1003.001 — OS Credential Dumping: LSASS Memory</mitre_technique>
  </question>

  <question id="3">
    <text>Which of the following BEST describes a LOLBin (Living Off the Land Binary)?</text>
    <options>
      <option id="A" correct="false">Malware that disguises itself as a system file by copying itself to the Windows folder</option>
      <option id="B" correct="true">A legitimate, signed system binary that attackers abuse to execute malicious actions</option>
      <option id="C" correct="false">A binary that only exists in memory and never touches disk</option>
      <option id="D" correct="false">An open-source offensive security tool bundled with Kali Linux</option>
    </options>
    <explanation>LOLBins (Living Off the Land Binaries) are legitimate, digitally signed operating system binaries that have capabilities (like downloading files, executing code, or making network connections) that attackers can repurpose for malicious activities. Examples include certutil.exe (download files), mshta.exe (execute remote scripts), and rundll32.exe (execute DLLs). Because they are legitimate system tools, they often bypass AV signature detection and application allowlisting.</explanation>
    <difficulty>basic</difficulty>
    <topic>Endpoint Attack Techniques</topic>
    <mitre_technique>T1218 — System Binary Proxy Execution</mitre_technique>
  </question>

  <question id="4">
    <text>In the Linux audit framework, the auid field has a value of 1000 in an event where uid=0 and euid=0. What does this tell a security analyst?</text>
    <options>
      <option id="A" correct="false">The event was generated by root performing a routine task</option>
      <option id="B" correct="true">The user account with UID 1000 ran sudo or su to gain root access; the original login session was from UID 1000</option>
      <option id="C" correct="false">The UID 0 and UID 1000 conflict indicates log corruption</option>
      <option id="D" correct="false">User 1000 has been granted permanent root access</option>
    </options>
    <explanation>The auid (audit UID) is set at login time and cannot be changed by the process itself — it tracks the original login identity. When auid=1000 but uid=0 and euid=0, it means: the session was started by user 1000 (login), who then used sudo or su to obtain root privileges. This is crucial for attribution — even if an attacker gains root, the auid field reveals which account they started from.</explanation>
    <difficulty>intermediate</difficulty>
    <topic>Linux Security Monitoring</topic>
    <mitre_technique>T1548.003 — Abuse Elevation Control Mechanism: Sudo and Sudo Caching</mitre_technique>
  </question>

  <question id="5">
    <text>An EDR platform generates an alert: w3wp.exe (IIS Worker Process) spawned cmd.exe. The command line is cmd.exe /c whoami. What type of attack does this MOST likely indicate?</text>
    <options>
      <option id="A" correct="false">An insider threat using the web server for administrative tasks</option>
      <option id="B" correct="true">A web application command injection or web shell exploitation</option>
      <option id="C" correct="false">A legitimate IIS health check running diagnostic commands</option>
      <option id="D" correct="false">A scheduled maintenance script configured in IIS</option>
    </options>
    <explanation>w3wp.exe (IIS Worker Process) is the ASP.NET application host. It should NEVER spawn cmd.exe or PowerShell in normal operation. When it does, it means user-supplied input (via HTTP request) was interpreted as an OS command — a classic command injection (CWE-78). The command whoami is the standard first command attackers run to confirm they have code execution and identify their privilege level. This is a near-certain true positive.</explanation>
    <difficulty>intermediate</difficulty>
    <topic>Web Application Security, EDR Detection</topic>
    <mitre_technique>T1190 — Exploit Public-Facing Application, T1059.003 — Windows Command Shell</mitre_technique>
  </question>

  <question id="6">
    <text>What is the key advantage of XDR (Extended Detection and Response) over traditional EDR?</text>
    <options>
      <option id="A" correct="false">XDR is less expensive than EDR</option>
      <option id="B" correct="false">XDR uses machine learning while EDR does not</option>
      <option id="C" correct="true">XDR correlates telemetry from multiple sources (endpoint, network, email, cloud) enabling cross-domain attack chain visibility</option>
      <option id="D" correct="false">XDR does not require agent deployment on endpoints</option>
    </options>
    <explanation>XDR extends EDR by ingesting and correlating telemetry from multiple security domains: endpoint (process, file, registry), network (firewall, proxy, NDR), email (phishing, attachment analysis), identity (AD, SSO), and cloud workloads. This enables analysts to see a complete attack chain — for example, linking a phishing email to the resulting process execution to the lateral movement to the data exfiltration — in a single incident view, rather than requiring analysts to manually correlate across separate tools.</explanation>
    <difficulty>basic</difficulty>
    <topic>EDR vs XDR Technology</topic>
  </question>

  <question id="7">
    <text>A threat hunter finds a Windows scheduled task named \Microsoft\Windows\EdgeUpdate\ScheduledUpdate that was created yesterday and runs PowerShell with a Base64-encoded command every 4 hours. The legitimate MicrosoftEdgeUpdate task has a different path and XML signature. Which MITRE ATT&CK technique is PRIMARILY being used?</text>
    <options>
      <option id="A" correct="false">T1059.001 — Command and Scripting Interpreter: PowerShell</option>
      <option id="B" correct="true">T1053.005 — Scheduled Task/Job: Scheduled Task combined with T1036.004 — Masquerade Task or Service</option>
      <option id="C" correct="false">T1105 — Ingress Tool Transfer</option>
      <option id="D" correct="false">T1190 — Exploit Public-Facing Application</option>
    </options>
    <explanation>While PowerShell (T1059.001) is also involved, the PRIMARY technique pattern here is Scheduled Task persistence (T1053.005) combined with Masquerading (T1036.004 — Masquerade Task or Service). The attacker created a scheduled task that mimics a legitimate Microsoft Edge Update task by using a similar name and path. T1036.004 specifically describes the technique of naming malicious tasks/services to look like legitimate ones. The PowerShell execution is the payload, but the primary detection-relevant technique is the persistence + masquerade combination.</explanation>
    <difficulty>advanced</difficulty>
    <topic>MITRE ATT&CK, Persistence Mechanisms</topic>
    <mitre_technique>T1053.005 — Scheduled Task, T1036.004 — Masquerade Task or Service</mitre_technique>
  </question>

</quiz>
